# bash-daemon/animation-daemon.service
[Unit]
Description=Steam Deck Custom Animation Daemon
After=network.target

[Service]
Type=simple
User=root
ExecStart=/home/deck/animation-daemon/animation-daemon.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target